# raku-sika-hub.service
[Unit]
Description=RakuSika Hub - Scale & Printer Gateway
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/raku-sika-hub
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=SCALE_DRIVER=mock
# Environment=SCALE_DRIVER=real
# Environment=SERIAL_PORT=/dev/ttyUSB0

[Install]
WantedBy=multi-user.target
