# deploy/raku-sika-hub.service
[Unit]
Description=RakuSika Hub - Scale & Printer Gateway
After=network.target avahi-daemon.service cups.service

[Service]
Type=simple
User=rakusika
WorkingDirectory=/home/rakusika/raku-sika-hub
ExecStart=/home/rakusika/raku-sika-hub/raku-sika-hub
Restart=always
RestartSec=5
Environment=HOME=/home/rakusika
# Environment=PRINTER_NAME=Brother_QL-820NWB
# Environment=MAX_CLIENTS=3
# Environment=SCALE_DRIVER=mock
# Environment=ALLOWED_ORIGINS=localhost:*,127.0.0.1:*,192.168.50.*
# Environment=ALLOW_ALL_ORIGINS=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raku-sika-hub

[Install]
WantedBy=multi-user.target
